# crontab
# shellcheck disable=all
PATH=/usr/local/bin:/usr/local/sbin:/usr/bin:/usr/sbin:/bin:/sbin
SHELL=/bin/bash

# 袩褉懈 褋褌邪褉褌械 褋懈褋褌械屑褘/泻芯薪褌械泄薪械褉邪
@reboot /bin/bash /app/startup.sh >> /var/log/cron.log 2>&1

# 袩芯写褗械屑 褉械蟹褞屑械 (褉邪蟹 胁 4 褔邪褋邪)
# 小懈屑胁芯谢 锌褉芯褑械薪褌邪 胁 泻褉芯薪褌邪斜械 薪褍卸薪芯 褝泻褉邪薪懈褉芯胁邪褌褜
# sleep 褋 褉邪薪写芯屑薪芯泄 蟹邪写械褉卸泻芯泄 芯褌 1 写芯 5 屑懈薪褍褌 薪褍卸械薪 胁芯-锌械褉胁褘褏 褔褌芯斜褘 褌芯泻械薪 芯斜薪芯胁懈褌褜褋褟 褍褋锌械谢, 邪 胁芯
# 胁褌芯褉褘褏, 褔褌芯斜褘 邪薪褌懈褋锌邪屑 薪械 褋褉邪斜芯褌邪谢
0 */5 * * * /bin/bash -c 'sleep $((60 + RANDOM \% 241))' && /usr/local/bin/python -m hh_applicant_tool update-resumes >> /var/log/cron.log 2>&1

# 袪邪褋褋褘谢邪械屑 芯褌泻谢懈泻懈 (褉械泻芯屑械薪写芯胁邪薪薪褘械 胁邪泻邪薪褋懈懈 锌芯 褉械蟹褞屑械, AI-褎懈谢褜褌褉)
0 8-21 * * * /bin/bash -c 'sleep $((60 + RANDOM \% 241))' && /usr/local/bin/python -m hh_applicant_tool apply-vacancies --letter /app/letter.txt -f --ai-filter heavy --excluded-filter 'junior|褋褌邪卸懈褉芯胁泻|bitrix|ddd|web3|crypto|blockchain|写褉褍卸薪\\w+泻芯谢谢械泻褌懈胁|锌芯谢懈谐褉邪褎|open\\s*space|芯锌械薪褋锌械泄褋|褏邪泻邪褌芯薪|泻芯薪泻褍褉褋|褌械褋褌芯胁\\w+ 蟹邪写邪薪|soft skill' >> /var/log/cron.log 2>&1

# 袪邪褋褋褘谢邪械屑 芯褌泻谢懈泻懈 锌芯 锌芯懈褋泻褍 "Python" 褋 AI-褎懈谢褜褌褉芯屑
30 8-21 * * * /bin/bash -c 'sleep $((60 + RANDOM \% 241))' && /usr/local/bin/python -m hh_applicant_tool apply-vacancies --search 'Python' --ai-filter heavy --letter /app/letter.txt -f --excluded-filter 'junior|褋褌邪卸懈褉芯胁泻|bitrix|ddd|web3|crypto|blockchain|写褉褍卸薪\\w+泻芯谢谢械泻褌懈胁|锌芯谢懈谐褉邪褎|open\\s*space|芯锌械薪褋锌械泄褋|褏邪泻邪褌芯薪|泻芯薪泻褍褉褋|褌械褋褌芯胁\\w+ 蟹邪写邪薪|soft skill' >> /var/log/cron.log 2>&1

# 孝械褋褌: 泻邪卸写褘械 5 屑懈薪褍褌 (蟹邪写械褉卸泻邪 写芯 1 屑懈薪褍褌褘, 褔褌芯斜褘 薪械 薪邪泻谢邪写褘胁邪谢懈褋褜)
#*/5 * * * * /bin/bash -c 'sleep $((1 + RANDOM \% 59))' && /usr/local/bin/python -m hh_applicant_tool apply-vacancies >> /var/log/cron.log 2>&1

# 袩褘褌邪械屑褋褟 芯斜薪芯胁懈褌褜 褌芯泻械薪 泻邪卸写褍褞 屑懈薪褍褌褍. 袙 褝褌芯屑 薪械褌 薪械芯斜褏芯写懈屑芯褋褌懈, 屑芯卸薪芯 蟹邪泻芯屑屑械薪褌懈褉芯胁邪褌褜, 褌邪泻 泻邪泻
# 褌芯泻械薪 芯斜薪芯胁谢褟械褌褋褟 邪胁褌芯屑邪褌懈褔械褋泻懈 锌褉懈 谢褞斜芯屑 蟹邪锌褉芯褋械
*/1 * * * * /usr/local/bin/python -m hh_applicant_tool refresh-token

# 袪邪褋泻芯屑屑械薪褌懈褉褍泄, 械褋谢懈 褏芯褔械褕褜 褍写邪谢褟褌褜 褔邪褌褘 褋 芯褌泻邪蟹邪屑懈
#0 */1 * * * /bin/bash -c 'sleep $((60 + RANDOM \% 241))' && /usr/local/bin/python -m hh-applicant-tool clear-negotiations -d >> /var/log/cron.log 2>&1

#*/1 * * * * echo "Cron is alive at $(date)" >> /var/log/cron.log 2>&1
